# templates/default/etc/profile.d/spark.sh.erb

PATH=$PATH:<%= @spark_path %>
export PATH

<% if ! @spark_public_dns.nil? -%>
SPARK_PUBLIC_DNS=<%= @spark_public_dns %>
export SPARK_PUBLIC_DNS
<% end -%>
